# cron/gzip_staticgen_cache.cron
*/17 * * * *  ~/project/cron/gzip_staticgen_cache.sh
